# homegear.conf
[program:homegear]
command=/usr/bin/homegear

